# crontab
0 0 * * *  root /usr/sbin/logrotate /etc/logrotate.d/suricata_rotate.conf
